# deploy/crontab
43 6 * * * /usr/bin/certbot renew --post-hook "systemctl reload nginx"
